# scripts/openclaw-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.openclaw.backup</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/jleechan/.openclaw/workspace/openclaw/scripts/run-openclaw-backup.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>14400</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/jleechan/Library/Logs/openclaw-backup/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jleechan/Library/Logs/openclaw-backup/stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/jleechan/.openclaw/workspace/openclaw</string>
  </dict>
</plist>
